# icecast2-liquidsoap/common/stubs/etc/systemd/system/icecast2.service
[Unit]
Description=Icecast2 daemon
Documentation=https://gitlab.xiph.org/xiph/icecast-server

[Service]
Type=simple
Restart=always
User=icecast
ExecStart=/usr/bin/icecast2 -b -c /etc/icecast2/icecast.xml
PIDFile=/run/icecast.pid
KillMode=process
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target